# config/cron/root.cron
# run this script on boot
@reboot /scripts/boot.sh
